# debian/prerm
#!/bin/bash
set -e

remove() {
    systemctl --user --global disable fpd-unlockd
    rm -f /etc/systemd/user/fpdlistener.service
}

case "$1" in
	"remove")
		remove
		;;
esac
